# var/lib/dpkg/info/apt-xapian-index.conffiles
/etc/cron.weekly/apt-xapian-index
/etc/dbus-1/system.d/org.debian.AptXapianIndex.conf
/etc/bash_completion.d/axi-cache
