# var/lib/dpkg/info/aptitude.list
/.
/var
/var/lib
/var/lib/aptitude
/usr
/usr/share
/usr/share/menu
/usr/share/menu/aptitude
/usr/share/bug
/usr/share/bug/aptitude
/usr/share/man
/usr/share/man/gl
/usr/share/man/gl/man8
/usr/share/man/gl/man8/aptitude-curses.8.gz
/usr/share/man/fi
/usr/share/man/fi/man8
/usr/share/man/fi/man8/aptitude-curses.8.gz
/usr/share/man/de
/usr/share/man/de/man8
/usr/share/man/de/man8/aptitude-curses.8.gz
/usr/share/man/ja
/usr/share/man/ja/man8
/usr/share/man/ja/man8/aptitude-curses.8.gz
/usr/share/man/man8
/usr/share/man/man8/aptitude-curses.8.gz
/usr/share/man/pl
/usr/share/man/pl/man8
/usr/share/man/pl/man8/aptitude-curses.8.gz
/usr/share/man/cs
/usr/share/man/cs/man8
/usr/share/man/cs/man8/aptitude-curses.8.gz
/usr/share/man/it
/usr/share/man/it/man8
/usr/share/man/it/man8/aptitude-curses.8.gz
/usr/share/man/fr
/usr/share/man/fr/man8
/usr/share/man/fr/man8/aptitude-curses.8.gz
/usr/share/man/es
/usr/share/man/es/man8
/usr/share/man/es/man8/aptitude-curses.8.gz
/usr/share/doc
/usr/share/doc/aptitude
/usr/share/doc/aptitude/changelog.gz
/usr/share/doc/aptitude/copyright
/usr/share/doc/aptitude/NEWS.Debian.gz
/usr/share/doc/aptitude/changelog.Debian.gz
/usr/bin
/usr/bin/aptitude-curses
/etc
/etc/cron.daily
/etc/cron.daily/aptitude
/etc/logrotate.d
/etc/logrotate.d/aptitude
