# var/lib/dpkg/info/avr-evtd.conffiles
/etc/avr-evtd/EventScript
/etc/avr-evtd/recovery.tar
/etc/avr-evtd/emergency-eth0
/etc/init.d/avr-evtd
/etc/default/avr-evtd
