# var/lib/dpkg/info/avr-evtd.list
/.
/usr
/usr/sbin
/usr/sbin/avr-evtd
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/avr-evtd.8.gz
/usr/share/doc
/usr/share/doc/avr-evtd
/usr/share/doc/avr-evtd/changelog.Debian.gz
/usr/share/doc/avr-evtd/NEWS.Debian.gz
/usr/share/doc/avr-evtd/copyright
/usr/share/doc/avr-evtd/README.IMPORTANT
/usr/share/doc/avr-evtd/changelog.gz
/etc
/etc/avr-evtd
/etc/avr-evtd/EventScript
/etc/avr-evtd/recovery.tar
/etc/avr-evtd/emergency-eth0
/etc/init.d
/etc/init.d/avr-evtd
/etc/default
/etc/default/avr-evtd
