# var/lib/dpkg/info/debsums.list
/.
/etc
/etc/cron.daily
/etc/cron.daily/debsums
/etc/cron.monthly
/etc/cron.monthly/debsums
/etc/default
/etc/default/debsums
/etc/debsums-ignore
/etc/cron.weekly
/etc/cron.weekly/debsums
/usr
/usr/sbin
/usr/sbin/debsums_init
/usr/bin
/usr/bin/rdebsums
/usr/bin/debsums
/usr/share
/usr/share/doc
/usr/share/doc/debsums
/usr/share/doc/debsums/README
/usr/share/doc/debsums/changelog.gz
/usr/share/doc/debsums/copyright
/usr/share/man
/usr/share/man/ru
/usr/share/man/ru/man8
/usr/share/man/ru/man8/debsums_init.8.gz
/usr/share/man/ru/man1
/usr/share/man/ru/man1/debsums.1.gz
/usr/share/man/fr
/usr/share/man/fr/man8
/usr/share/man/fr/man8/debsums_init.8.gz
/usr/share/man/fr/man1
/usr/share/man/fr/man1/debsums.1.gz
/usr/share/man/de
/usr/share/man/de/man8
/usr/share/man/de/man8/debsums_init.8.gz
/usr/share/man/de/man1
/usr/share/man/de/man1/debsums.1.gz
/usr/share/man/sv
/usr/share/man/sv/man8
/usr/share/man/sv/man8/debsums_init.8.gz
/usr/share/man/sv/man1
/usr/share/man/sv/man1/debsums.1.gz
/usr/share/man/man8
/usr/share/man/man8/debsums_init.8.gz
/usr/share/man/pt
/usr/share/man/pt/man8
/usr/share/man/pt/man8/debsums_init.8.gz
/usr/share/man/pt/man1
/usr/share/man/pt/man1/debsums.1.gz
/usr/share/man/man1
/usr/share/man/man1/rdebsums.1.gz
/usr/share/man/man1/debsums.1.gz
/usr/share/man/es
/usr/share/man/es/man8
/usr/share/man/es/man8/debsums_init.8.gz
/usr/share/man/es/man1
/usr/share/man/es/man1/debsums.1.gz
