# var/lib/dpkg/info/debsums.conffiles
/etc/cron.daily/debsums
/etc/cron.monthly/debsums
/etc/default/debsums
/etc/debsums-ignore
/etc/cron.weekly/debsums
