# var/lib/dpkg/info/dpkg.conffiles
/etc/logrotate.d/dpkg
/etc/alternatives/README
/etc/dpkg/dpkg.cfg
/etc/cron.daily/dpkg
