# var/lib/dpkg/info/ifupdown.conffiles
/etc/init.d/ifupdown
/etc/init.d/ifupdown-clean
/etc/default/ifupdown
