# var/lib/dpkg/info/initscripts.list
/.
/bin
/bin/mountpoint
/sys
/sbin
/sbin/fsck.nfs
/usr
/usr/share
/usr/share/initscripts
/usr/share/initscripts/default.rcS
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/halt.5.gz
/usr/share/man/man5/rcS.5.gz
/usr/share/man/man1
/usr/share/man/man1/mountpoint.1.gz
/usr/share/man/man8
/usr/share/man/man8/fsck.nfs.8.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/initscripts
/usr/share/doc
/usr/share/doc/initscripts
/usr/share/doc/initscripts/changelog.gz
/usr/share/doc/initscripts/copyright
/usr/share/doc/initscripts/README.Debian
/usr/share/doc/initscripts/changelog.Debian.gz
/etc
/etc/network
/etc/network/if-up.d
/etc/network/if-up.d/mountnfs
/etc/default
/etc/default/bootlogd
/etc/default/tmpfs
/etc/default/halt
/etc/default/devpts
/etc/init.d
/etc/init.d/mountall-bootclean.sh
/etc/init.d/bootlogd
/etc/init.d/single
/etc/init.d/skeleton
/etc/init.d/mountdevsubfs.sh
/etc/init.d/sendsigs
/etc/init.d/bootmisc.sh
/etc/init.d/rc.local
/etc/init.d/bootlogs
/etc/init.d/killprocs
/etc/init.d/umountroot
/etc/init.d/stop-bootlogd
/etc/init.d/checkfs.sh
/etc/init.d/checkroot.sh
/etc/init.d/mountnfs.sh
/etc/init.d/rmnologin
/etc/init.d/hostname.sh
/etc/init.d/umountnfs.sh
/etc/init.d/mountkernfs.sh
/etc/init.d/mountnfs-bootclean.sh
/etc/init.d/mountoverflowtmp
/etc/init.d/urandom
/etc/init.d/halt
/etc/init.d/mountall.sh
/etc/init.d/umountfs
/etc/init.d/mtab.sh
/etc/init.d/reboot
/etc/init.d/stop-bootlogd-single
/lib
/lib/init
/lib/init/swap-functions.sh
/lib/init/vars.sh
/lib/init/usplash-fsck-functions.sh
/lib/init/rw
/lib/init/mount-functions.sh
/lib/init/splash-functions-base
/lib/init/bootclean.sh
/var
/var/lib
/var/lib/initscripts
/var/lib/urandom
/var/log
/var/log/fsck
