# var/lib/dpkg/info/libdevmapper1.02.list
/etc
/etc/init.d
/etc/init.d/libdevmapper1.02
