# var/lib/dpkg/info/linuxlogo.conffiles
/etc/init.d/linuxlogo
/etc/linux_logo.conf
