# var/lib/dpkg/info/linuxlogo.list
/.
/usr
/usr/bin
/usr/bin/linux_logo
/usr/share
/usr/share/locale
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pt_BR/LC_MESSAGES/linux_logo.mo
/usr/share/locale/pl
/usr/share/locale/pl/LC_MESSAGES
/usr/share/locale/pl/LC_MESSAGES/linux_logo.mo
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/it/LC_MESSAGES/linux_logo.mo
/usr/share/locale/uk
/usr/share/locale/uk/LC_MESSAGES
/usr/share/locale/uk/LC_MESSAGES/linux_logo.mo
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/sv/LC_MESSAGES/linux_logo.mo
/usr/share/locale/nl
/usr/share/locale/nl/LC_MESSAGES
/usr/share/locale/nl/LC_MESSAGES/linux_logo.mo
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/linux_logo.mo
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/es/LC_MESSAGES/linux_logo.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/linux_logo.mo
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/linux_logo.mo
/usr/share/locale/zh_TW
/usr/share/locale/zh_TW/LC_MESSAGES
/usr/share/locale/zh_TW/LC_MESSAGES/linux_logo.mo
/usr/share/locale/da
/usr/share/locale/da/LC_MESSAGES
/usr/share/locale/da/LC_MESSAGES/linux_logo.mo
/usr/share/locale/nb
/usr/share/locale/nb/LC_MESSAGES
/usr/share/locale/nb/LC_MESSAGES/linux_logo.mo
/usr/share/pixmaps
/usr/share/pixmaps/linuxlogo.xpm
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/linux_logo.1.gz
/usr/share/menu
/usr/share/menu/linuxlogo
/usr/share/doc
/usr/share/doc/linuxlogo
/usr/share/doc/linuxlogo/USAGE.gz
/usr/share/doc/linuxlogo/changelog.gz
/usr/share/doc/linuxlogo/copyright
/usr/share/doc/linuxlogo/LINUX_LOGO.FAQ.gz
/usr/share/doc/linuxlogo/examples
/usr/share/doc/linuxlogo/examples/login_script
/usr/share/doc/linuxlogo/examples/logo2ascii.pl
/usr/share/doc/linuxlogo/examples/named_pipe
/usr/share/doc/linuxlogo/examples/mklogo
/usr/share/doc/linuxlogo/CHANGES_IN_5.0
/usr/share/doc/linuxlogo/README.CUSTOM_LOGOS.gz
/usr/share/doc/linuxlogo/README.Debian
/usr/share/doc/linuxlogo/README.SECURITY
/usr/share/doc/linuxlogo/TODO
/usr/share/doc/linuxlogo/ANNOUNCE.logo
/usr/share/doc/linuxlogo/BUGS
/usr/share/doc/linuxlogo/changelog.Debian.gz
/usr/share/doc/linuxlogo/README.gz
/etc
/etc/linux_logo.conf
/etc/init.d
/etc/init.d/linuxlogo
/usr/bin/linuxlogo
/usr/share/man/man1/linuxlogo.1.gz
