# var/lib/dpkg/info/minissdpd.list
/.
/usr
/usr/sbin
/usr/sbin/minissdpd
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/minissdpd.1.gz
/usr/share/doc
/usr/share/doc/minissdpd
/usr/share/doc/minissdpd/changelog.Debian.gz
/usr/share/doc/minissdpd/copyright
/usr/share/doc/minissdpd/README
/usr/share/doc/minissdpd/changelog.gz
/etc
/etc/init.d
/etc/init.d/minissdpd
/etc/default
/etc/default/minissdpd
