# var/lib/dpkg/info/minissdpd.prerm
#!/bin/sh
set -e
# Automatically added by dh_installinit
if [ -x "/etc/init.d/minissdpd" ]; then
	invoke-rc.d minissdpd stop || exit $?
fi
# End automatically added section
