# var/lib/dpkg/info/nfs-kernel-server.list
/.
/var
/var/lib
/var/lib/nfs
/var/lib/nfs/v4recovery
/usr
/usr/share
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/exports.5.gz
/usr/share/man/man7
/usr/share/man/man7/nfsd.7.gz
/usr/share/man/man8
/usr/share/man/man8/nfsd.8.gz
/usr/share/man/man8/exportfs.8.gz
/usr/share/man/man8/mountd.8.gz
/usr/share/man/man8/svcgssd.8.gz
/usr/share/doc
/usr/share/doc/nfs-kernel-server
/usr/share/doc/nfs-kernel-server/changelog.gz
/usr/share/doc/nfs-kernel-server/copyright
/usr/share/doc/nfs-kernel-server/NEWS.Debian.gz
/usr/share/doc/nfs-kernel-server/changelog.Debian.gz
/usr/share/nfs-kernel-server
/usr/share/nfs-kernel-server/conffiles
/usr/share/nfs-kernel-server/conffiles/nfs-kernel-server.default.md5sum
/usr/share/nfs-kernel-server/conffiles/etc.exports
/usr/share/nfs-kernel-server/conffiles/nfs-kernel-server.default
/usr/share/nfs-kernel-server/conffiles/etc.exports.md5sum
/usr/sbin
/usr/sbin/rpc.svcgssd
/usr/sbin/rpc.nfsd
/usr/sbin/exportfs
/usr/sbin/rpc.mountd
/etc
/etc/init.d
/etc/init.d/nfs-kernel-server
/usr/share/man/man8/rpc.nfsd.8.gz
/usr/share/man/man8/rpc.svcgssd.8.gz
/usr/share/man/man8/rpc.mountd.8.gz
