# var/lib/dpkg/info/openbsd-inetd.prerm
#!/bin/sh
set -e
# Automatically added by dh_installinit
if [ -x "/etc/init.d/openbsd-inetd" ]; then
	invoke-rc.d openbsd-inetd stop || exit $?
fi
# End automatically added section
