# var/lib/dpkg/info/samba.conffiles
/etc/cron.daily/samba
/etc/network/if-up.d/samba
/etc/init.d/samba
/etc/logrotate.d/samba
