# var/lib/dpkg/info/transmission-daemon.list
/.
/var
/var/lib
/var/lib/transmission-daemon
/var/lib/transmission-daemon/downloads
/var/lib/transmission-daemon/info
/usr
/usr/share
/usr/share/bug
/usr/share/bug/transmission-daemon
/usr/share/bug/transmission-daemon/presubj
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/transmission-daemon.1.gz
/usr/share/doc
/usr/bin
/usr/bin/transmission-daemon
/etc
/etc/transmission-daemon
/etc/transmission-daemon/settings.json
/etc/transmission-daemon/README.json
/etc/init.d
/etc/init.d/transmission-daemon
/etc/default
/etc/default/transmission-daemon
/var/lib/transmission-daemon/info/settings.json
/usr/share/doc/transmission-daemon
