# var/lib/dpkg/info/transmission-daemon.conffiles
/etc/transmission-daemon/settings.json
/etc/transmission-daemon/README.json
/etc/init.d/transmission-daemon
/etc/default/transmission-daemon
