# var/lib/dpkg/info/util-linux.list
/.
/bin
/bin/tailf
/bin/more
/bin/dmesg
/sbin
/sbin/hwclock
/sbin/mkfs.minix
/sbin/raw
/sbin/mkfs.bfs
/sbin/getty
/sbin/switch_root
/sbin/agetty
/sbin/mkswap
/sbin/ddisk
/sbin/cfdisk
/sbin/sfdisk
/sbin/mkfs.cramfs
/sbin/blockdev
/sbin/fsck
/sbin/findfs
/sbin/blkid
/sbin/fsck.minix
/sbin/isosize
/sbin/mkfs
/sbin/wipefs
/sbin/fsck.cramfs
/sbin/ctrlaltdel
/sbin/pivot_root
/usr
/usr/bin
/usr/bin/ipcmk
/usr/bin/mcookie
/usr/bin/pg
/usr/bin/chrt
/usr/bin/lscpu
/usr/bin/setarch
/usr/bin/unshare
/usr/bin/taskset
/usr/bin/delpart
/usr/bin/setterm
/usr/bin/line
/usr/bin/fallocate
/usr/bin/namei
/usr/bin/flock
/usr/bin/chkdupexe
/usr/bin/getopt
/usr/bin/rev
/usr/bin/ipcrm
/usr/bin/ipcs
/usr/bin/ddate
/usr/bin/rename.ul
/usr/bin/whereis
/usr/bin/addpart
/usr/bin/ionice
/usr/bin/setsid
/usr/bin/partx
/usr/share
/usr/share/info
/usr/share/info/ipc.info.gz
/usr/share/locale
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pl
/usr/share/locale/pl/LC_MESSAGES
/usr/share/locale/pl/LC_MESSAGES/util-linux-ng.mo
/usr/share/locale/eu
/usr/share/locale/eu/LC_MESSAGES
/usr/share/locale/eu/LC_MESSAGES/util-linux-ng.mo
/usr/share/locale/hu
/usr/share/locale/hu/LC_MESSAGES
/usr/share/locale/hu/LC_MESSAGES/util-linux-ng.mo
/usr/share/locale/fi
/usr/share/locale/fi/LC_MESSAGES
/usr/share/locale/ca
/usr/share/locale/ca/LC_MESSAGES
/usr/share/locale/zh_CN
/usr/share/locale/zh_CN/LC_MESSAGES
/usr/share/locale/zh_CN/LC_MESSAGES/util-linux-ng.mo
/usr/share/locale/tr
/usr/share/locale/tr/LC_MESSAGES
/usr/share/locale/vi
/usr/share/locale/vi/LC_MESSAGES
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/id
/usr/share/locale/id/LC_MESSAGES
/usr/share/locale/et
/usr/share/locale/et/LC_MESSAGES
/usr/share/locale/ja
/usr/share/locale/ja/LC_MESSAGES
/usr/share/locale/uk
/usr/share/locale/uk/LC_MESSAGES
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/nl
/usr/share/locale/nl/LC_MESSAGES
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/da
/usr/share/locale/da/LC_MESSAGES
/usr/share/locale/cs
/usr/share/locale/cs/LC_MESSAGES
/usr/share/locale/sl
/usr/share/locale/sl/LC_MESSAGES
/usr/share/man
/usr/share/man/man5
/usr/share/man/man1
/usr/share/man/man1/fallocate.1.gz
/usr/share/man/man1/namei.1.gz
/usr/share/man/man1/linux64.1.gz
/usr/share/man/man1/unshare.1.gz
/usr/share/man/man1/taskset.1.gz
/usr/share/man/man1/readprofile.1.gz
/usr/share/man/man1/line.1.gz
/usr/share/man/man1/ionice.1.gz
/usr/share/man/man1/setsid.1.gz
/usr/share/man/man1/rename.ul.1.gz
/usr/share/man/man1/more.1.gz
/usr/share/man/man1/dmesg.1.gz
/usr/share/man/man1/ipcmk.1.gz
/usr/share/man/man1/whereis.1.gz
/usr/share/man/man1/getopt.1.gz
/usr/share/man/man1/lscpu.1.gz
/usr/share/man/man1/setterm.1.gz
/usr/share/man/man1/chrt.1.gz
/usr/share/man/man1/ddate.1.gz
/usr/share/man/man1/mcookie.1.gz
/usr/share/man/man1/rev.1.gz
/usr/share/man/man1/flock.1.gz
/usr/share/man/man1/ipcs.1.gz
/usr/share/man/man1/tailf.1.gz
/usr/share/man/man1/pg.1.gz
/usr/share/man/man1/linux32.1.gz
/usr/share/man/man1/chkdupexe.1.gz
/usr/share/man/man1/ipcrm.1.gz
/usr/share/man/man8
/usr/share/man/man8/ldattach.8.gz
/usr/share/man/man8/rdev.8.gz
/usr/share/man/man8/hwclock.8.gz
/usr/share/man/man8/setarch.8.gz
/usr/share/man/man8/raw.8.gz
/usr/share/man/man8/ddisk.8.gz
/usr/share/man/man8/ramsize.8.gz
/usr/share/man/man8/fsck.minix.8.gz
/usr/share/man/man8/rootflags.8.gz
/usr/share/man/man8/isosize.8.gz
/usr/share/man/man8/addpart.8.gz
/usr/share/man/man8/vidmode.8.gz
/usr/share/man/man8/ctrlaltdel.8.gz
/usr/share/man/man8/findfs.8.gz
/usr/share/man/man8/pivot_root.8.gz
/usr/share/man/man8/cfdisk.8.gz
/usr/share/man/man8/mkfs.minix.8.gz
/usr/share/man/man8/agetty.8.gz
/usr/share/man/man8/switch_root.8.gz
/usr/share/man/man8/fdformat.8.gz
/usr/share/man/man8/wipefs.8.gz
/usr/share/man/man8/mkswap.8.gz
/usr/share/man/man8/ppc.8.gz
/usr/share/man/man8/delpart.8.gz
/usr/share/man/man8/blockdev.8.gz
/usr/share/man/man8/fsck.8.gz
/usr/share/man/man8/tunelp.8.gz
/usr/share/man/man8/mkfs.bfs.8.gz
/usr/share/man/man8/cytune.8.gz
/usr/share/man/man8/ppc64.8.gz
/usr/share/man/man8/getty.8.gz
/usr/share/man/man8/blkid.8.gz
/usr/share/man/man8/rtcwake.8.gz
/usr/share/man/man8/ppc32.8.gz
/usr/share/man/man8/sfdisk.8.gz
/usr/share/man/man8/mkfs.8.gz
/usr/share/man/man8/partx.8.gz
/usr/share/man/man3
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/util-linux
/usr/share/doc
/usr/share/doc/util-linux
/usr/share/doc/util-linux/changelog.gz
/usr/share/doc/util-linux/copyright
/usr/share/doc/util-linux/examples
/usr/share/doc/util-linux/examples/getopt-parse.tcsh
/usr/share/doc/util-linux/examples/sfdisk.examples.gz
/usr/share/doc/util-linux/examples/getopt-test.bash
/usr/share/doc/util-linux/examples/getopt-parse.bash
/usr/share/doc/util-linux/examples/fstab.example2
/usr/share/doc/util-linux/examples/getopt-test.tcsh
/usr/share/doc/util-linux/README.getty
/usr/share/doc/util-linux/README.modems-with-agetty
/usr/share/doc/util-linux/README.Debian.hwclock
/usr/share/doc/util-linux/README.fdisk.gz
/usr/share/doc/util-linux/changelog.Debian.gz
/usr/share/doc/util-linux/README.poeigl.gz
/usr/share/doc/util-linux/AUTHORS.gz
/usr/share/doc/util-linux/README.cfdisk
/usr/sbin
/usr/sbin/fdformat
/usr/sbin/rtcwake
/usr/sbin/readprofile
/usr/sbin/rdev
/usr/sbin/tunelp
/usr/sbin/ldattach
/usr/sbin/cytune
/usr/lib
/usr/lib/mime
/usr/lib/mime/packages
/usr/lib/mime/packages/util-linux
/etc
/etc/init.d
/etc/init.d/hwclock.sh
/etc/init.d/hwclockfirst.sh
/lib
/lib/udev
/lib/udev/hwclock-set
/lib/udev/rules.d
/lib/udev/rules.d/85-hwclock.rules
/usr/bin/ppc64
/usr/bin/linux32
/usr/bin/ppc32
/usr/bin/ppc
/usr/bin/linux64
/usr/sbin/vidmode
/usr/sbin/ramsize
/usr/sbin/rootflags
