# var/lib/dpkg/info/vsftpd.list
/.
/etc
/etc/ftpusers
/etc/pam.d
/etc/pam.d/vsftpd
/etc/logrotate.d
/etc/logrotate.d/vsftpd
/etc/vsftpd.conf
/etc/init.d
/etc/init.d/vsftpd
/usr
/usr/share
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/vsftpd.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/vsftpd.8.gz
/usr/share/doc
/usr/share/doc/vsftpd
/usr/share/doc/vsftpd/REWARD
/usr/share/doc/vsftpd/README.security
/usr/share/doc/vsftpd/TUNING
/usr/share/doc/vsftpd/README.ssl
/usr/share/doc/vsftpd/SIZE
/usr/share/doc/vsftpd/FAQ.gz
/usr/share/doc/vsftpd/README
/usr/share/doc/vsftpd/changelog.Debian.gz
/usr/share/doc/vsftpd/copyright
/usr/share/doc/vsftpd/changelog.gz
/usr/share/doc/vsftpd/README.Debian
/usr/share/doc/vsftpd/SPEED
/usr/share/doc/vsftpd/examples
/usr/share/doc/vsftpd/examples/INTERNET_SITE_NOINETD
/usr/share/doc/vsftpd/examples/INTERNET_SITE_NOINETD/README
/usr/share/doc/vsftpd/examples/INTERNET_SITE_NOINETD/vsftpd.conf
/usr/share/doc/vsftpd/examples/VIRTUAL_USERS_2
/usr/share/doc/vsftpd/examples/VIRTUAL_USERS_2/README
/usr/share/doc/vsftpd/examples/README
/usr/share/doc/vsftpd/examples/INTERNET_SITE
/usr/share/doc/vsftpd/examples/INTERNET_SITE/README.gz
/usr/share/doc/vsftpd/examples/INTERNET_SITE/vsftpd.conf
/usr/share/doc/vsftpd/examples/INTERNET_SITE/vsftpd.xinetd
/usr/share/doc/vsftpd/examples/VIRTUAL_USERS
/usr/share/doc/vsftpd/examples/VIRTUAL_USERS/logins.txt
/usr/share/doc/vsftpd/examples/VIRTUAL_USERS/README.gz
/usr/share/doc/vsftpd/examples/VIRTUAL_USERS/vsftpd.pam
/usr/share/doc/vsftpd/examples/VIRTUAL_USERS/vsftpd.conf
/usr/share/doc/vsftpd/examples/vsftpd.preseed
/usr/share/doc/vsftpd/examples/PER_IP_CONFIG
/usr/share/doc/vsftpd/examples/PER_IP_CONFIG/hosts.allow
/usr/share/doc/vsftpd/examples/PER_IP_CONFIG/README
/usr/share/doc/vsftpd/examples/VIRTUAL_HOSTS
/usr/share/doc/vsftpd/examples/VIRTUAL_HOSTS/README
/usr/share/doc/vsftpd/SECURITY
/usr/share/doc/vsftpd/SECURITY/IMPLEMENTATION
/usr/share/doc/vsftpd/SECURITY/TRUST.gz
/usr/share/doc/vsftpd/SECURITY/DESIGN.gz
/usr/share/doc/vsftpd/SECURITY/OVERVIEW
/usr/share/doc/vsftpd/AUDIT
/usr/share/doc/vsftpd/BUGS
/usr/share/doc/vsftpd/BENCHMARKS
/usr/share/doc/vsftpd/TODO
/usr/share/doc/vsftpd/NEWS.Debian.gz
/usr/sbin
/usr/sbin/vsftpd
